# .mounts/home-sieben-share.automount
[Unit]
Description=Automount ZMBFS:/tank/share

[Automount]
Where=/home/sieben/share
TimeoutIdleSec=30

[Install]
WantedBy=multi-user.target