# user/mountfiles/latest-mountfiles/syno-backup.automount
[Unit]
Description=Automount syno-backup
Wants=network-online.target
After=network-online.target

[Automount]
Where=/syno/backup

[Install]
WantedBy=multi-user.target
