# user/archive/mountfiles-v2/syno-archiv.automount
[Unit]
Description=Automount syno-archiv

[Automount]
Where=/syno/archiv

[Install]
WantedBy=multi-user.target
