# user/archive/mountfiles-v3/syno-archiv.automount
[Unit]
Description=Automount syno-archiv
Wants=network-online.target
After=network-online.target

[Automount]
Where=/syno/archiv

[Install]
WantedBy=multi-user.target
