# user/archive/mountfiles-v3/syno-music.automount
[Unit]
Description=Automount syno-music
Wants=network-online.target
After=network-online.target

[Automount]
Where=/syno/music

[Install]
WantedBy=multi-user.target
