# user/archive/mountfiles-v3/syno-cloudsync.automount
[Unit]
Description=Automount syno-cloudsync
Wants=network-online.target
After=network-online.target

[Automount]
Where=/syno/cloudsync

[Install]
WantedBy=multi-user.target
