# user/archive/mountfiles-v3/syno-downloads.automount
[Unit]
Description=Automount syno-downloads
Wants=network-online.target
After=network-online.target

[Automount]
Where=/syno/downloads

[Install]
WantedBy=multi-user.target
