# etc/systemd/system/rsync-backup-sieben.timer
[Unit]
Description=Run Rsync Backup Sieben Service daily

[Timer]
OnBootSec=5min 
Persistent=true

[Install]
WantedBy=timers.target
