# etc/systemd/system/syno-home.automount
[Unit]
Description=Automount syno-home
Wants=network-online.target
After=network-online.target

[Automount]
Where=/syno/home

[Install]
WantedBy=multi-user.target
