# etc/systemd/system/syno-home.timer
[Unit]
Description=Timer mount syno-home

[Timer]
OnBootSec=10s
Unit=syno-home.mount

[Install]
WantedBy=multi-user.target
